# files/etc/uci-defaults/network
# This file is to get rid of the annoying network-default
# configuration!!!
